# com.congressai.refresh.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.congressai.refresh</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/python3.11</string>
        <string>/Users/dylanmoskowitz/Downloads/congress_ai/refresh.py</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTORCH_ENABLE_MPS_FALLBACK</key>
        <string>1</string>
        <key>PYTORCH_MPS_HIGH_WATERMARK_RATIO</key>
        <string>0.0</string>
    </dict>

    <!-- Run every day at 6:00 AM -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>6</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/dylanmoskowitz/Downloads/congress_ai/logs/refresh.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/dylanmoskowitz/Downloads/congress_ai/logs/refresh_error.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
